# crontab
*/10	* * * *	wget -nv --spider localhost:8000/cache
0		* * * *	wget -nv --spider localhost:8000/recommendations
